# etc/systemd/ltsm_connector.socket
[Socket]
ListenStream=0.0.0.0:3389
ListenStream=0.0.0.0:5900
Accept=yes

[Install]
WantedBy=sockets.target
