# etc/systemd/ltsm_connector_rdp.socket
[Socket]
ListenStream=3389
Accept=yes

[Install]
WantedBy=sockets.target
