# etc/systemd/ltsm_service.service
[Unit]
Description=LTSM service

[Service]
ExecStart=/usr/local/sbin/LTSM_service

[Install]
WantedBy=multi-user.target
